# supervisord.ini
[supervisord]
nodaemon=true

[program:conf-watch]
command=/opt/conf-watch/confd-watch -h etcd
autorestart=true
stdout_logfile=/logs/confd-watch.log

[program:nginx]
command=/usr/sbin/nginx
autorestart=true
